# data/linux_tasks/systemd_user/alpha-pilot-doctor.timer
[Unit]
Description=AlphaPilot Doctor每5分钟巡检

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Unit=alpha-pilot-doctor.service
Persistent=true

[Install]
WantedBy=timers.target
