# data/linux_tasks/systemd_user/quant-pilot-doctor.service
[Unit]
Description=Quant Pilot Watchdog巡检/自愈/闭环修复
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/env bash /home/ubuntu/projects/quant-pilot/data/linux_tasks/run_doctor.sh
Environment=BROKER_MODE=paper
Environment=PYTHONUNBUFFERED=1
