# data/linux_tasks/systemd_user/quant-pilot-auto-restart.service
[Unit]
Description=Quant Pilot 自动盘Watchdog重启
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/env bash /home/ubuntu/projects/quant-pilot/data/linux_tasks/restart_auto.sh
Environment=BROKER_MODE=paper
Environment=PYTHONUNBUFFERED=1
