# scripts/quant-realtime.service
[Unit]
Description=AlphaPilot Realtime Trader
After=network.target

[Service]
Type=simple
User=%u
WorkingDirectory=%h/projects/alpha-pilot
Environment="PATH=%h/projects/alpha-pilot/.venv/bin:%h/projects/alpha-pilot/venv/bin:%h/.local/bin:/usr/local/bin:/usr/bin:/bin"
Environment=BROKER_MODE=paper
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=%h/.hermes/.env
ExecStart=/usr/bin/env python3 main.py --realtime
Restart=always
RestartSec=10
StandardOutput=append:%h/.hermes/logs/quant-realtime.log
StandardError=append:%h/.hermes/logs/quant-realtime-error.log

[Install]
WantedBy=default.target
